# connot.service
[Unit]
Description=ConnNotify notifier - forwards collector events to KDE
After=graphical-session.target
Wants=dbus.socket

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /opt/connot/connot_notifier.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=default.target
